# FocusTimer.app/Contents/Resources/com.focustimer.helper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.focustimer.helper</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/juns/focus_timer/focus_timer_env/bin/python</string>
        <string>/Users/juns/focus_timer/FocusTimer.app/Contents/MacOS/FocusTimerHelper</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/var/log/FocusTimer/helper.log</string>

    <key>StandardErrorPath</key>
    <string>/var/log/FocusTimer/helper_error.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/juns/focus_timer</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>ProcessType</key>
    <string>Background</string>

    <key>ThrottleInterval</key>
    <integer>60</integer>

    <key>ExitTimeOut</key>
    <integer>10</integer>
</dict>
</plist>